# php-fpm/7.1/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/bin/php-fpm -F
user=root
